# contrib/bytobot.supervisor.conf
[program:bytebot]
command=/path/to/bytebot/env/bin/python /path/to/bytebot/bytebot.py
user=bytebot

